# deploy/pincerna.service
[Unit]
Description=Pincerna Flask API
After=network.target

[Service]
User=root
Group=www-data
WorkingDirectory=/opt/pincerna/services/api
Environment="PATH=/opt/pincerna/.venv/bin"
EnvironmentFile=-/opt/pincerna/.env
ExecStart=/opt/pincerna/.venv/bin/gunicorn \
    --workers 2 \
    --bind 127.0.0.1:5002 \
    --timeout 120 \
    --access-logfile /var/log/pincerna/access.log \
    --error-logfile /var/log/pincerna/error.log \
    app:app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target